# trunk/pathfinder-debian/debian/pathfinderd.prerm
#!/bin/sh

if [ -x "/etc/init.d/pathfinderd" ]; then
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d pathfinderd stop || exit $?
	else
		/etc/init.d/pathfinderd stop || exit $?
	fi
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#
